# reprodebian/mkosi.skeleton/etc/systemd/system/nvidia-persistenced.service.d/10-cmdline.conf
[Unit]
ConditionKernelCommandLine=cocoon_gpu
